# make_MagOS/files/patches/rootfs/rootfs-MagOS-2012/lib/systemd/system/nfs.target
[Unit]
Description=Network File System Server
Requires=var-lib-nfs-rpc_pipefs.mount pcbind.service
After=network.target named.service 

[Install]
WantedBy=multi-user.target
