# make_MagOS/files/patches/rootfs/rootfs/usr/share/magos/initrd/etc/profile.d/20screen.sh
if [ -z "$SCREENDIR" ]; then
    export SCREENDIR=$HOME/tmp
fi
