# make_MagOS/files/patches/rootfs/rootfs/usr/lib/magos/rc.local.d/81-samba
#!/bin/bash

### BEGIN INIT INFO
# Provides: rc.local.MagOS.samba
# X-Mandriva-Compat-Mode
# Default-Start: 2 3 4 5
# Short-Description: Generate initial passwords for samba users
# Description: This script will generate initial passwords for samba users
### END INIT INFO

#. /etc/rc.d/init.d/functions

case "$1" in
  start)
     if ! grep -q nobody /etc/samba/passdb.tdb ;then
        echo -e \\n\\n | smbpasswd -Lsa nobody
        for a in `grep users /etc/group | awk -F: '{print $4}' | tr ';,' ' '` ;do
           PASS=`echo $a | cat - /mnt/live/VERSION /proc/cmdline | md5sum | cut -c 1-5`
           echo -e $PASS\\n$PASS\\n | smbpasswd -Lsa $a >/dev/null 2>&1
        done
     fi
     ;;
  showdefpass)
        for a in `grep users /etc/group | awk -F: '{print $4}' | tr ';,' ' '` ;do
           PASS=`echo $a | cat - /mnt/live/VERSION /proc/cmdline | md5sum | cut -c 1-5`
           echo $a $PASS
        done
     ;;
esac
