# make_MagOS/files/patches/rootfs/rootfs/etc/profile.d/numlock.sh
# Linux-Mandrake configuration.
#	ReLock the NumLock key if /etc/rc.d/init.d/numlock has been run.
#	This is needed because login resets the tty.

