# make_modules/add-modules/noupstream_patches/rootfs/etc/cron.d/ntpdate
# Sync time every 10 minutes if enable
*/10 * * * *  root /usr/lib/magos/cron/ntpdate.sh
