# make_modules/source/multimedia/xbmc/mountemall.rules
SUBSYSTEM=="block", KERNEL=="sd*", ACTION=="add", ENV{ID_FS_USAGE}=="filesystem|other",  RUN+="/usr/lib/magos/scripts/pmount-hal-user /dev/%k"
